# build/template/runner/unregister.socket
[Unit]
Description=Privilege isolation intermediary between control HTTP server and the OS

[Install]
WantedBy=sockets.target

[Socket]
ListenFIFO=/run/unregister-runners.fifo
